# vagrant/scionupgrade.timer
[Unit]
Description=Run SCION upgrade daily

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=multi-user.target
